# conf/crontab
* * * * * root cp -r /var/www/vhosts/* /var/backup/
